# robophery/files/supervisor.conf
{%- from "robophery/map.jinja" import server with context %}

[program:robophery_manager]
directory={{ server.dir.base }}
command={{ server.dir.base }}/bin/rp_manager
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=syslog
stderr_logfile=syslog
